# src/main/scripts/init.d_opsagent
#!/bin/bash
#
# /etc/init.d/opsagent
#
# Startup script for Ops agent
#
# chkconfig: 2345 20 80
# description: Starts and stops Ops agent

. /etc/rc.d/init.d/functions

export OPSAGENT_HOME=/opt/opsagent
export OPSAGENT_CONF=/opt/opsagent/conf
export OPSAGENT_OWNR=cassandra

log_file=/var/log/opsagent/opsagent.log
pid_file=/var/run/opsagent.pid
OPSAGENT_PROG=/opt/opsagent/bin/opsagent

case "$1" in
    start)
        # Ops agent startup
        touch $pid_file
        chown -R $OPSAGENT_OWNR.$OPSAGENT_OWNR $pid_file
        echo -n "Starting Ops agent: "
        su $OPSAGENT_OWNR -c "$OPSAGENT_PROG -p $pid_file" > $log_file 2>&1
        echo "OK"
        ;;
    stop)
        # Ops agent shutdown
        echo -n "Shutdown Ops agent: "
        su $OPSAGENT_OWNR -c "kill `cat $pid_file`"
        #su $OPSAGENT_OWNR -c "kill -9 `ps -ef | grep opsagent |grep -v grep | awk '{print $2}'`"
        #for t in `seq 40`; do $0 status > /dev/null 2>&1 && sleep 0.5 || break; done
        #for t in `seq 40`; do netstat -tnlp | grep "0.0.0.0:8080" > /dev/null 2>&1 && sleep 0.1 || break; done
        sleep 5
        echo "OK"
        ;;
    reload|restart)
        $0 stop
        $0 start
        ;;
    status)
        PIDS=$(ps -ef | grep opsagent | grep -v grep | awk '{print $2}')
        ps -p $PIDS
        exit $?
        ;;
    *)
        echo "Usage: `basename $0` start|stop|status|restart|reload"
        exit 1
esac

exit 0